# docker/php/supervisord.example.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=storage/app/supervisord.pid

[program:cli]
# workerman https://manual.workerman.net/doc/zh-cn/install/start-and-stop.html
command=php app/Workerman/http.php start
# append the `-d` when production

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=5

[include]
files = /etc/supervisor/conf.d/*.conf